# crontab/crontjob.e
# 必須寫絕對路徑，一定要指定使用的 command

*/1 * * * * /bin/sh /opt/echo_date.sh
