# crontab/echo_date.sh
#!/bin/sh
now="$(date +'%Y%m%d')"
echo "$now"
# ...:.20200303

# date +"%Y-%m-%d %T"
# 2019-03-25 14:40:32

# ....
echo $now >> /opt/now.txt
